# deploy/systemd/floonet-authority.service
# Hardened systemd unit for the Floonet name authority on bare metal.
#
# Recommended: the guided installer wires all of this up for you, including the
# secret credential drop-in below:
#   sudo ./install.sh            # builds, installs this unit, runs the wizard
#
# By hand:
#   cd name-authority && cargo build --release
#   sudo install -m0755 target/release/floonet-name-authority /usr/local/bin/
#   sudo floonet-name-authority setup        # writes /etc/floonet-authority.env
#   sudo install -m0644 ../deploy/systemd/floonet-authority.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now floonet-authority
#
# The service stores only public data plus payment grant state, but it is
# still locked down: dynamic unprivileged user, read-only system, no new
# privileges. A free authority (FLOONET_PAY_MODE=off) holds NO secret at all.
#
# When a paid mode is configured, the GoblinPay token is a secret and must never
# sit in the (root:root 0640, but proxy/backup-exposable) env file. The wizard
# writes it to /etc/floonet-authority/secrets/goblinpay_token at mode 0600 and
# the installer drops in floonet-authority.service.d/10-goblinpay-credential.conf
# with:
#   [Service]
#   LoadCredential=goblinpay_token:/etc/floonet-authority/secrets/goblinpay_token
#   Environment=GOBLINPAY_TOKEN_FILE=%d/goblinpay_token
# so systemd (as root) loads the 0600 file and exposes a copy under
# $CREDENTIALS_DIRECTORY (%d) that the dynamic service user can read. That
# Environment= overrides the /etc path the env file names. The base unit stays
# valid for a free deploy because the LoadCredential lives only in the drop-in.

[Unit]
Description=Floonet name authority (name@domain -> nostr pubkey)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
# DynamicUser allocates a throwaway unprivileged user at runtime. If you need
# a stable owner for the data dir, comment this out and set `User=floonet`
# (create the user first).
DynamicUser=yes

# Identity/config. Edit /etc/floonet-authority.env (copy of .env.example).
EnvironmentFile=/etc/floonet-authority.env

# Managed state at /var/lib/floonet-authority (created and chowned by systemd).
StateDirectory=floonet-authority
StateDirectoryMode=0750
Environment=FLOONET_NAMES_DB=/var/lib/floonet-authority/names.db

ExecStart=/usr/local/bin/floonet-name-authority
Restart=on-failure
RestartSec=2

# --- hardening ---
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
# Only the state directory is writable.
ReadWritePaths=/var/lib/floonet-authority
# No raw sockets; only IP.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
