# deploy/us-east/floonet-authority.service.d/10-us-east.conf
# us-east production overrides for the generic hardened unit
# (deploy/systemd/floonet-authority.service).
#
# The box keeps every Goblin service's data under /opt/goblin (a single backup
# root) and runs them as the unprivileged `goblin` account, so we swap the
# generic unit's DynamicUser for a stable owner and relocate the writable path.
# The base unit's own comment sanctions exactly this ("If you need a stable
# owner for the data dir ... set User="). Every other hardening directive from
# the base unit (ProtectSystem=strict, NoNewPrivileges, ProtectHome, the
# @system-service syscall filter, RestrictAddressFamilies, ...) is inherited
# unchanged — goblin's home is /opt/goblin/nip05d, not under /home, so
# ProtectHome=yes stays safe.
[Service]
DynamicUser=no
User=goblin
Group=goblin

# The base unit ships StateDirectory=floonet-authority plus a matching
# /var/lib writable path. Clear both (empty assignment resets the list) and
# point everything at the /opt/goblin tree instead.
StateDirectory=
ReadWritePaths=
ReadWritePaths=/opt/goblin/floonet-authority
WorkingDirectory=/opt/goblin/floonet-authority

# Applied after the base unit, so this wins over the base Environment= line.
Environment=FLOONET_NAMES_DB=/opt/goblin/floonet-authority/names.db
